# rtsvd.ini
; Sample supervisor config file.

[unix_http_server]
file=./supervisor.sock   ; (the path to the socket file)

[supervisord]
logfile=./supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=30MB       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=3          ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=./supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=true               ; (start in foreground if true;default false)
directory=./
childlogdir=./ ; ('AUTO' child log dir, default $TEMP)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix://./supervisor.sock ; use a unix:// URL  for a unix socket

[program:roundtable]
command=python3 -m wxagent.mainrt start irc
autostart=true
autorestart=false
redirect_stderr=true

[program:irc]
command=python3 -m wxagent.mainrt start roundtable
autostart=true
autorestart=false
redirect_stderr=true

[program:tox]
command=python3 -m wxagent.mainrt start tox
autostart=true
autorestart=false
redirect_stderr=true

[program:tox]
command=python3 -m wxagent.mainrt start wechat
autostart=true
autorestart=false
redirect_stderr=true
